# bundles/sentry/sentry.cron
#!/bin/sh

su sentry -s /bin/sh -c 'sentry cleanup --days=30'
